# deploy/coinex-bot.service
[Unit]
Description=CoinEX_SkamerList Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/coinex-bot
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
